# deploy/tgmail-relay.service
[Unit]
Description=Telegram IMAP Relay (real-time via IDLE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/tgmail-relay
ExecStart=/usr/local/bin/tgmail-relay
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
